# src/management/cron/src/crontabs/root
# min   hour    day   month   weekday   command
0       0       *     *       *         echo "Should attempt renew"
*/5     *       *     *       *         /app/scripts/request-cert.sh
